# install/etc/udev/rules.d/99-input.rules
SUBSYSTEM=="input", RUN+="/usr/lib/udev/bluetooth"
